# systemd-unit/deltabot.service
[Unit]
Description=Deltachat-Bot-Gatekeeper Service
After=network.target bluetooth.service

[Service]
Type=simple
# Runs as an unprivileged user. BLE access relies on file capabilities
# on bluepy-helper:
#   sudo setcap cap_net_raw,cap_net_admin+eip \
#     <venv>/lib/python3.11/site-packages/bluepy/bluepy-helper
# See README section 8 ("Running as an unprivileged user").
User=pi
Group=pi
ExecStart=/home/pi/gatekeeper-bot/start-gatekeeper-bot.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
